# scripts/profile.d/webpage.sh
# enter proper PATH to the webpage root directory (where repo was cloned)
export WEBPAGE_PATH=
